# ansible/roles/smarthome/files/crontab
1-59/2 * * * * root /usr/bin/python /var/lib/smarthome/modbus2.py >> /var/lib/smarthome/modbus2.log
